# xdg/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Type=Application
Exec=flameshot
Comment=Screenshot tool
X-GNOME-Autostart-enabled=true
